# source/etc/supervisor.conf
[program:art]
command=
  /srv/art/env/bin/uwsgi
  /srv/art/etc/uwsgi.ini
autorestart=true
stopsignal=INT
