# {{cookiecutter.project_name}}/etc/supervisor.conf
[program:art]
command=
  /srv/art/env/bin/uwsgi
  /srv/art/source/etc/uwsgi.ini
autorestart=true
stopsignal=INT
